# deploy/pcb-inspect.service
[Unit]
Description=PCB Inspect — Zone Check Application
After=network.target

[Service]
User=pcbinspect
Group=pcbinspect
WorkingDirectory=/opt/pcb-inspect
EnvironmentFile=/opt/pcb-inspect/.env
ExecStart=/opt/pcb-inspect/venv/bin/gunicorn -c gunicorn.conf.py app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/pcb-inspect/uploads
PrivateTmp=true

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
